# claude-code-ui.service.template
[Unit]
Description=Claude Code UI Server
After=network.target

[Service]
Type=simple
WorkingDirectory=$PWD
Environment=PATH=$HOME/.nvm/versions/node/$NODE_VERSION/bin:/usr/local/bin:/usr/bin:/bin
Environment=PORT=3001
Environment=NODE_ENV=production
ExecStart=$HOME/.nvm/versions/node/$NODE_VERSION/bin/npm run start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
